# scripts/Xavier/crontab
* * * * * /root/checkCopy.sh >/dev/null 2>&1 #CHANGEME
* * * * * sleep 30;/root/checkCopy.sh >/dev/null 2>&1 #CHANGEME