# scripts/initScripts/dOCR.service
[Unit]
Description=OCR service
After=overclock.service
[Service]
Type=simple
ExecStart=/root/initScripts/dOCR.sh #CHANGEME
