# scripts/initScripts/dVehiclePerson.service
[Unit]
Description=Vehicle Person service
After=overclock.service
[Service]
Type=simple
ExecStart=/root/initScripts/dVehiclePerson.sh #CHANGEME
