# scripts/initScripts/overclock.service
[Unit]
Description=Overclock service
After=systemd-user-sessions.service
[Service]
Type=simple
ExecStart=/root/initScripts/overclock.sh #CHANGEME
