# scripts/initScripts/watcher.service
[Unit]
Description=Watcher service
After=watcher.service
[Service]
Type=simple
ExecStart=/root/initScripts/watcher.sh #CHANGEME
